# scripts/imaging-servo.cron
@reboot root /bin/bash -c "python /usr/libexec/imaging-servo/ntp-update"
