# scripts/imaging-worker.cron
@reboot root /bin/bash -c "python /usr/libexec/eucalyptus-imaging-worker/dns-update"
